# ops/gpu-metrics.timer
[Unit]
Description=Run GPU metrics exporter every 15s

[Timer]
OnBootSec=10
OnUnitActiveSec=15s
AccuracySec=1s

[Install]
WantedBy=timers.target
