# ops/server.service
[Unit]
Description=Orion HAL HTTP server
After=network-online.target vllm.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/orion
ExecStart=%h/orion/.venv/bin/python -m hal.server
EnvironmentFile=%h/orion/.env
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
