# ops/telegram.service
# ROLLBACK ONLY — not the primary deployment path.
# The Telegram bot runs inside the Docker container alongside the HTTP
# server, managed by supervisord.
# This unit file is retained for emergency rollback to host-venv mode.
# See OPERATIONS.md for the current deployment procedure.

[Unit]
Description=Orion Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/orion
ExecStart=%h/orion/.venv/bin/python -m hal.telegram
EnvironmentFile=%h/orion/.env
Restart=on-failure
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
