# ops/vllm.service
[Unit]
Description=vLLM — Qwen2.5-32B-Instruct-AWQ OpenAI-compatible API
After=network.target

[Service]
Type=simple
User=jp
WorkingDirectory=%h
Environment=VLLM_USE_FLASHINFER_SAMPLER=0
ExecStart=%h/vllm-env/bin/vllm serve Qwen/Qwen2.5-32B-Instruct-AWQ \
    --enable-auto-tool-choice \
    --tool-call-parser hermes \
    --port 8000 \
    --enforce-eager \
    --max-model-len 8192 \
    --gpu-memory-utilization 0.95
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
