# apps/rag/unified/scripts/orion-update.timer
[Unit]
Description=ORION Weekly Dataset Update Timer
Documentation=https://github.com/your-repo/orion-rag
Requires=orion-update.service

[Timer]
# Run every Sunday at 2:00 AM
OnCalendar=Sun *-*-* 02:00:00

# If the system was powered off during the scheduled time, run on next boot
Persistent=true

# Random delay up to 30 minutes to avoid simultaneous updates
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
